# src/vodafone-mobile-connect/contrib/option_icon/files/99-option-icon.rules
# Simone Tolotti
SYSFS{idVendor}=="05c6", SYSFS{idProduct}=="1000", RUN+="/usr/sbin/icon_switch"

