# attic/VodafoneMobileConnectCard/build/opt/vmc/etc/udev/rules.d/99-option-icon.rules
# Simone Tolotti
SYSFS{idVendor}=="05c6", SYSFS{idProduct}=="1000", RUN+="/opt/vmc/usr/sbin/icon_switch"

