# src/VodafoneMobileConnectCard/output/src/opt/vmc/etc/udev/rules.d/99-novatel-ovation.rules
SUBSYSTEMS=="scsi", KERNEL=="sr[0-9]*", ATTRS{vendor}=="Novatel ", SYMLINK="novatel-storage"
SYSFS{idVendor}=="1410", SYSFS{idProduct}=="5010", RUN+="/opt/vmc/usr/sbin/ovation.sh"
