# src/vodafone-mobile-connect/resources/platform/eeepc/etc/udev/rules.d/45-vmc-option.rules
SYSFS{idVendor}=="05c6", SYSFS{idProduct}=="1000", RUN+="/usr/sbin/vmc-option.sh 05c6 1000"

